# twitter-dashboard.service
[Unit]
Description=AI-TwitterPersona Web Dashboard
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=twitterbot
Group=twitterbot
WorkingDirectory=/opt/ai-twitterpersona
Environment=PATH=/opt/ai-twitterpersona/venv/bin
ExecStart=/opt/ai-twitterpersona/venv/bin/gunicorn --config gunicorn.conf.py app:app --workers=1
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/ai-twitterpersona

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-twitterpersona-dashboard

[Install]
WantedBy=multi-user.target